# debian_files/kvkafkalistener.upstart
author "Vegard Bones <vegard.bones@met.no>"

start on stopped rc
stop on runlevel [016]

respawn
respawn limit 10 60

# For some reason a straighforward exec with setuid does not work on newly 
# created machines before a reboot has been performed 
exec start-stop-daemon --start -c kvalobs --exec /usr/bin/kvkafkalistener -- -o /var/log/kvalobs/kafkastatus.txt
