# src/tools/nagios/kvkafkalistener.in.service
[Unit]
Description=Runs the kvkafaklistener

[Service]
Type=simple
ExecStart=__BINPATH__/kvkafkalistener -o /var/log/kvalobs/kafkastatus.txt
